# files/etc/init.d/fastdick
#!/bin/sh /etc/rc.common

# Writen by zz090923610
# Not for any kind of commercial use.

START=99
STOP=10
SERVICE_WRITE_PID=1
SERVICE_DAEMONIZE=1

start()
{
	service_start /usr/bin/fastdick
}

stop()
{
	cd /proc
	for pid in [0-9]*
	do
		if grep "Name:	fastdick" $pid/status &>/dev/null
			then
			if [ "$pid" -eq "$$" ]
				then
				continue
			fi
			echo "Stop service fastdick $pid"
			kill "$pid"
			for cpid in [0-9]*
			do
				if grep "PPid:	$pid" $cpid/status &>/dev/null
				then
				    echo "Stop child process $cpid"
				    name=`cat $cpid/status | grep "Name:" | sed 's/Name://' | tr -d ' \n\t'`
				    if [ "$name" = "sleep" ]
				    then
					    kill $cpid
				    fi
				fi
			done	
		fi
	done	
}

restart()
{
	stop
	start
}
